# deployment/salt/ss_wave/files/ss_wave.service
[Unit]
Description=Seascape Wave Agent
Before=multi-user.target

[Service]
Type=simple
ExecStart=/var/run/ss_wave.env/bin/python -m ss_wave
User=root
Group=root
WorkingDirectory=/opt/seascape/seascape_wave

Environment="SS_WEBSOCKET_URL={{ salt.pillar.get('ss_wave:websocket_url', 'wss://ingest.seascape.example/ingest/websocket') }}"
Environment="SS_API_KEY={{ salt.pillar.get('ss_wave:api_key', '00000000-0000-0000-0000-000000000000') }}"
Environment="SS_INSECURE={{ salt.pillar.get('ss_wave:insecure', 'True') }}"

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
